# bin/screen.crontab
#min    hour    mday    month   wday    command
*/5 * * * * /Users/andersonfreitas/bin/screen.sh 

